# scripts/bootGovernor.service
[Unit]
Description=bootGovernor

[Service]
ExecStart=/usr/bin/bootGovernor

[Install]
WantedBy=multi-user.target
